# nuc/monitoring/presencia-monitor.service
[Unit]
Description=Presencia Health Monitor Report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/presencia/monitoring/presencia-monitor.sh
